# swupdate.service
[Unit]
Description=SWUpdate daemon

[Service]
ExecStart=/usr/sbin/swupdate.sh
KeyringMode=inherit
